# scripts/tendermate-frontend.service
[Unit]
Description=TenderMate Frontend (Vite dev server)
After=network.target

[Service]
Type=simple
User=krzysztof-balicki
WorkingDirectory=/home/krzysztof-balicki/Nextcloud/GIT/przetargi/frontend
ExecStart=/usr/bin/npx vite --host --port 5173
Restart=always
RestartSec=3
Environment=NODE_ENV=development

[Install]
WantedBy=default.target
